# templates/nginx/botblocker/ngxblocker-update.service.epp
[Unit]
Description=Update the ngxblocker data

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-ngxblocker -q
User=root
